# profile.d/dircolors.sh
eval $( dircolors -b /etc/profile.d/dircolors )
