# profile.d/dircolors.sh

eval $( dircolors -b $DEVTOOLS/ali.md/.dircolors )
