# profile.d/logo.sh
if [[ -f /etc/profile.d/logo.txt ]]
then
  cat /etc/profile.d/logo.txt
fi
